# autostart/mainline-notify.desktop
[Desktop Entry]
Exec=bash "/home/gustavo/.config/mainline/mainline-notify.sh" --autostart
